# launchd/com.awvm.reaper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.awvm.reaper</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/uv</string>
    <string>run</string>
    <string>/Users/pst/code/pst/pstaylor-patrick/azure-windows/scripts/awvm.py</string>
    <string>reaper</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/pst/code/pst/pstaylor-patrick/azure-windows</string>
  <key>StartInterval</key>
  <integer>900</integer>
  <key>StandardOutPath</key>
  <string>/Users/pst/code/pst/pstaylor-patrick/azure-windows/.azure-windows/reaper.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/pst/code/pst/pstaylor-patrick/azure-windows/.azure-windows/reaper.log</string>
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
